# sortd.service
[Unit]
Description=sortd context-aware file organizer daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/sortd daemon start
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
